# deploy/wisp-mount-share.service
# Persistent shared bind mount for wisp's FUSE propagation.
#
# Problem this solves: wisp serves /mnt/wisp via FUSE. Media-server containers
# (Silo, Plex, …) bind that path in with :rslave so they see wisp's mount. But if
# the shared peer group lives on the FUSE mount itself, every wisp restart (and
# every host reboot) destroys it — orphaning the consumers, which then see
# "Transport endpoint is not connected" until THEY restart too.
#
# The fix: create a durable self-bind at /mnt/wisp *before* Docker starts, and
# make it rshared. The peer group then lives on this bind, not on wisp's FUSE
# mount, so it survives wisp restarts and reboots. wisp mounts its FUSE on top;
# consumers slave to the durable peer group and keep seeing the current mount.
#
# Install (stop wisp first so /mnt/wisp is empty when the durable bind is created —
# the guard below skips the bind if the path is already a FUSE mountpoint):
#   docker compose stop wisp
#   sudo cp deploy/wisp-mount-share.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now wisp-mount-share.service
#   docker compose up -d wisp
# Then recreate the media-server containers once so they slave to the new group.
# (At boot the unit runs before Docker, so /mnt/wisp is empty and the bind is made
# correctly without any manual stop.)

[Unit]
Description=Persistent shared bind mount for wisp FUSE propagation
Documentation=https://github.com/dreulavelle/wisp
DefaultDependencies=no
After=local-fs.target
Before=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/mkdir -p /mnt/wisp
# Bind only if it isn't already a mountpoint (idempotent across restarts).
ExecStart=/bin/sh -c 'mountpoint -q /mnt/wisp || mount --bind /mnt/wisp /mnt/wisp'
ExecStart=/bin/mount --make-rshared /mnt/wisp
# Never tear the peer group down while consumers depend on it.
ExecStop=/bin/true

[Install]
WantedBy=multi-user.target
